# service/dsync-daemon.service
[Unit]
Description=Dsync Daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/dsync-daemon
Restart=always

[Install]
WantedBy=multi-user.target